# profile/airootfs/etc/systemd/system/installer.service
[Unit]
Description=Tateru Installer

[Service]
Type=simple
ExecStart=/etc/rc.installer

[Install]
WantedBy=multi-user.target
